# systemd/netnat.service
[Unit]
Description=NetNAT - Proxmox NAT & Port Forwarding Manager
Documentation=https://github.com/yourusername/netnat
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/netnat
ExecStart=/opt/netnat/netnat -config /etc/netnat/config.yml
Restart=always
RestartSec=5
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netnat

[Install]
WantedBy=multi-user.target